# sys-block/targetcli-fb/files/target.init
#!/sbin/openrc-run
command=/usr/bin/targetcli
name=targetcli
description="iSCSI target configuration persistence"  
extra_started_commands="save"

depend() 
{
  true
}
  
start() 
{
  /usr/bin/targetcli restoreconfig
}
  
stop() 
{
  /usr/bin/targetcli clearconfig confirm=True
}

save()
{
  /usr/bin/targetcli saveconfig
}